# dotfiles/.config/systemd/system/disable-usb-wakeup.service
[Unit]
Description=Configure USB Wakeup for Sleep Stability
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable-usb-wakeup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
